# packaging/centos7/startup/dialer.service
[Unit]
Description=Job that runs the dialer daemon

[Service]
Type=simple
ExecStart=/usr/bin/dialer run ahn

[Install]
WantedBy=multi-user.target
